# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:apache]
command=apachectl start
autorestart=true

[program:mysql]
command=mysqld_safe
autorestart=true
